# systemd/redshift-day.service
[Unit]
Description=Redshift Day Mode

[Service]
Type=oneshot
ExecStart=/usr/bin/redshift -x -m randr
